# etc/systemd/system/umount@.service
[Unit]

[Service]
Type=oneshot
ExecStartPre=/bin/umount "/run/media/%I"
ExecStart=/bin/rm -r "/run/media/%I"
